# system_health/crontabs/crontab.redis
* * * * * root bash -c "/system_health/scheduler.sh redis" >> /var/log/cron_redis.log 2>&1
# empty line