# system_health/crontabs/crontab.es
* * * * * root bash -c "/system_health/scheduler.sh es" >> /var/log/cron_es.log 2>&1
# empty line