# system_health/crontabs/crontab.celery
* * * * * root bash -c "/system_health/scheduler.sh celery" >> /var/log/cron_celery.log 2>&1
# empty line